# .bin/com.hash-thing.crew-nudge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hash-thing.crew-nudge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/edward/projects/hash-thing/.bin/crew-nudge</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/edward/projects/hash-thing</string>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/edward/projects/hash-thing/.ship-notes/crew-nudge-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/edward/projects/hash-thing/.ship-notes/crew-nudge-launchd.log</string>
</dict>
</plist>
